# crontab
*/15 * * * * python /app/launch.py >> /var/log/cron.log 2>&1
